# php_monolith/src/Compass/Pivot/sh/cron/crontab.cron
# general
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/general.php start bot0 > /dev/null 2>&1 &

# phphooker [0-9]
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/phphooker.php start bot0 > /dev/null 2>&1 &
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/phphooker.php start bot1 > /dev/null 2>&1 &
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/phphooker.php start bot2 > /dev/null 2>&1 &
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/phphooker.php start bot3 > /dev/null 2>&1 &
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/phphooker.php start bot4 > /dev/null 2>&1 &
* * * * *  /usr/local/bin/php /app/src/Compass/Pivot/api/cron/phphooker.php start bot5 > /dev/null 2>&1 &

# cron_sms_dispatcher
* * * * *  /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/sms/dispatcher.php start bot0 now > /dev/null 2>&1 &

# cron_mail_dispatcher
* * * * *  /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/mail/dispatcher.php start bot0 now > /dev/null 2>&1 &

# cron_sms_provider_observer
* * * * *  /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/sms/provider/observer.php start bot0 now > /dev/null 2>&1 &


# cron company observer
* * * * * /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/company/observer.php start bot0 now > /dev/null 2>&1 &

# cron company task
* * * * * /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/company/task.php start task-type-exit now > /dev/null 2>&1 &

# cron company tariff task
* * * * * /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/company/tariff_plan_task.php start bot0 now > /dev/null 2>&1 &

# cron service company task
* * * * * /usr/local/bin/php /app/src/Compass/Pivot/api/cron/company/service_task.php start bot0 > /dev/null 2>&1 &
* * * * * /usr/local/bin/php /app/src/Compass/Pivot/api/cron/company/service_task.php start bot1 > /dev/null 2>&1 &
* * * * * /usr/local/bin/php /app/src/Compass/Pivot/api/cron/company/service_task.php start bot2 > /dev/null 2>&1 &

# cron_system_reserve_server_company_observe
* * * * * /usr/local/bin/php -dxdebug.start_with_request=0 /app/src/Compass/Pivot/api/cron/system/reserve_server_company_observe.php start bot0 now > /dev/null 2>&1 &
