# php_monolith/src/Compass/Jitsi/sh/cron/crontab.cron
# general
* * * * *  /usr/local/bin/php /app/src/Compass/Jitsi/api/cron/general.php start bot0 > /dev/null 2>&1 &
